# packages/bluesky-notify/bluesky_notify-0.5.2-py3-none-any.whl/bluesky_notify/services/systemd/bluesky-notify.service
[Unit]
Description=Bluesky Notify Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/bluesky-notify start
Restart=always
RestartSec=10
User=%i

[Install]
WantedBy=default.target
